# debian-php-83/supervisor/services/php-fpm.conf
[program:php-fpm]
command = %(ENV_PHP_FPM_BIN)s -F -R -y %(ENV_PHP_FPM_CONF)s
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
